# root/etc/s6-overlay/s6-rc.d/env/up
/command/with-contenv /bin/bash -c /app/prepare-nginx.sh
