# root/etc/s6-overlay/s6-rc.d/uid/up
/command/with-contenv /bin/bash -c /app/setuid.sh