# root/etc/s6-overlay/s6-rc.d/backend/run
#!/command/with-contenv bash
cd /data
exec node index.js
